# crontab
MAILTO=root

01 * * * * root @libexecdir@/php-session.sh
