# crontab
MAILTO=root

01 * * * * root /usr/lib/php-session.sh
